# agent/com.tombridger.work-live.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tombridger.work-live</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>/Users/tombridger/.config/work-live/run-capture.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <!-- Every 5 minutes by wall clock, NOT StartInterval. Per launchd.plist(5),
       an interval firing that lands while the system is asleep "will be
       missed due to shortcomings in kqueue(3)" - observed as a permanently
       pended job after multi-day sleeps (no ticks at all the morning back).
       Calendar intervals instead coalesce missed firings into one run on
       wake, so the first wake after 8am always produces the morning
       check-in. Cadence stays 5 min = lib/time.ts captureIntervalMinutes. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Minute</key><integer>0</integer></dict>
    <dict><key>Minute</key><integer>5</integer></dict>
    <dict><key>Minute</key><integer>10</integer></dict>
    <dict><key>Minute</key><integer>15</integer></dict>
    <dict><key>Minute</key><integer>20</integer></dict>
    <dict><key>Minute</key><integer>25</integer></dict>
    <dict><key>Minute</key><integer>30</integer></dict>
    <dict><key>Minute</key><integer>35</integer></dict>
    <dict><key>Minute</key><integer>40</integer></dict>
    <dict><key>Minute</key><integer>45</integer></dict>
    <dict><key>Minute</key><integer>50</integer></dict>
    <dict><key>Minute</key><integer>55</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/work-live-agent.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/work-live-agent.err.log</string>
</dict>
</plist>
